# sts-robot/rc.local


# for autostart of the remote control add this to rc.local
cd /home/pi/robot
python remote-control.py &
