# simbamon/simbamon
#!/usr/bin/env bash
#
# simbamon -- a simple battery monitor
#
# Author: Hamish Cunningham <hamish@gate.ac.uk>
# This code is copyright Hamish Cunningham and the University of Sheffield
# and is licenced under GPL 3 or any later version.
#
### BEGIN INIT INFO
# Provides:          simbamon
# Short-Description: SimBaMon: a simple battery monitor daemon
# Description:       This script is implements the daemon for
#                    SimBaMon (a simple battery monitor)
### END INIT INFO

# Defaults
PATH=/sbin:/usr/sbin:/bin:/usr/bin
DESC="SimBaMon: a simple battery monitor daemon"
NAME=simbamon
DAEMON=/usr/sbin/$NAME
DAEMON_ARGS=""
#DAEMON_ARGS="--options args"
PIDFILE=/var/run/$NAME.pid
SCRIPTNAME=/etc/init.d/$NAME

# Read configuration variable file if it is present
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

# Pretty date/time function
pdate() { date +%b-%d-%Y-%T; }

# Tell the world
logger "$0: running at `pdate`"
logger "${NAME}: monitor frequency is ${MONITOR_FREQUENCY} seconds"

# Length of time to wait between routine log messages
i=0
LOG_INTERVAL=20
# TODO for production should be much higher

# Current level; -1 means unset
BAT_LEVEL=-1

# Do the work
while :
do
  # check the level
  # TODO

  # routine log messages
  i=`expr ${i} + 1`
  [ ${i} -eq ${LOG_INTERVAL} ] && i=0 && \
    logger "${NAME}: battery level is ${BAT_LEVEL} at `pdate`"

  # take a break
  sleep ${MONITOR_FREQUENCY}
done
